# etc/supervisor/bootstrap.ini
[program:bootstrap]
user=aska
command=/home/aska/.rbenv/shims/unicorn /home/aska/7kai-Bootstrap/config.ru -c /home/aska/7kai-Bootstrap/unicorn.conf.rb -E production -p 10001
# command=foreman start
directory=/home/aska/7kai-Bootstrap
stdout_logfile=/var/log/supervisor/bootstrap.log
stderr_logfile=/var/log/supervisor/bootstrap-error.log
numprocs = 1
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile_backups = 5
stderr_logfile_backups = 5
autostart = true
autorestart = true
stopsignal = QUIT
startsecs = 5
priority = 998